# debian/cron.daily
#!/bin/sh
django-admin cleanup --settings=www.settings --pythonpath=/var/www/django/songscape/
